# config/upstart/storm-ui.conf
start on runlevel [2345]
stop on runlevel [!2345]
 
env JAVA_HOME=/usr/lib/jvm/java-6-sun

script 
  cd /usr/local/storm/
  /usr/local/storm/bin/storm ui
end script

respawn
